# deploy/systemd/alpaca-telegram-integrity.service
[Unit]
Description=Alpaca Telegram + data integrity cycle (milestone 250, coverage, strict gate)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/stock-bot
Environment=TRADING_BOT_ROOT=/root/stock-bot
EnvironmentFile=-/root/stock-bot/.env
EnvironmentFile=-/root/.alpaca_env
ExecStart=/root/stock-bot/venv/bin/python3 /root/stock-bot/scripts/run_alpaca_telegram_integrity_cycle.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
