# deploy/systemd/paper-exec-mode-worker.timer
[Unit]
Description=Run paper-exec-mode-worker every minute during experiments

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=30s

[Install]
WantedBy=timers.target
