# deploy/systemd/telegram-failure-detector.timer
[Unit]
Description=Run Telegram failure detector every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true
Unit=telegram-failure-detector.service

[Install]
WantedBy=timers.target
