# deploy/systemd/telegram-failure-detector.service
[Unit]
Description=Telegram failure pager (Alpaca post-close / direction / integrity expectations)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/stock-bot
Environment=TRADING_BOT_ROOT=/root/stock-bot
EnvironmentFile=-/root/stock-bot/.env
# Block pager Telegram in prod; alpaca-telegram-integrity.service does not set this (allowlisted script_names still send).
Environment=TELEGRAM_GOVERNANCE_INTEGRITY_ONLY=1
ExecStart=/root/stock-bot/venv/bin/python3 /root/stock-bot/scripts/governance/telegram_failure_detector.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
