# deploy/systemd/uw-flow-daemon.service
[#] uw-flow-daemon.service
[#] Systemd unit to guarantee single-instance UW polling.

[Unit]
Description=Unusual Whales Flow Daemon (single instance)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/stock-bot
EnvironmentFile=-/root/stock-bot/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/stock-bot/venv/bin/python /root/stock-bot/uw_flow_daemon.py
Restart=always
RestartSec=5

# Safety: ensure a stuck process can be stopped
TimeoutStopSec=20
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target

