# deploy/systemd/wheel-broker-reconcile.timer
[Unit]
Description=Run wheel broker reconcile every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
AccuracySec=1min
Unit=wheel-broker-reconcile.service

[Install]
WantedBy=timers.target
