# deploy/stock-bot-dashboard.service
# stock-bot-dashboard.service
# Dashboard (Flask on port 5000). Restart after deploy to load new routes.
# Install: sudo cp /root/stock-bot/deploy/stock-bot-dashboard.service /etc/systemd/system/
#         sudo systemctl daemon-reload && sudo systemctl enable stock-bot-dashboard && sudo systemctl start stock-bot-dashboard

[Unit]
Description=STOCK-BOT Dashboard (Flask :5000)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/stock-bot
EnvironmentFile=-/root/stock-bot/.env
Environment=PYTHONUNBUFFERED=1
Environment=PORT=5000
ExecStart=/usr/bin/python3 /root/stock-bot/dashboard.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
